# etc/init.d/lightdm-wallpaper
#!/sbin/openrc-run
# Copyright 2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="lightdm-wallpaper"
description="Change le fond d'écran que lightdm prendra"

depend() {
	need localmount 
}
start() {
	einfo "Setting up lightdm wallpaper"
	files=(/home/woomy/Pictures/wallpapers/*)
	num=${#files[@]}
	wall="${files[RANDOM % num]}"
	if [ -f "/home/woomy/lightdm-wallpaper.jpg" ]; then
		rm "/home/woomy/lightdm-wallpaper.jpg"
	fi
	cp "${wall}" "/home/woomy/lightdm-wallpaper.jpg"
}
stop() {
	einfo "This service can't be stopped"
}
